# roles/appserver/files/upstart/uwsgi-emperor.conf.j2
start on net-device-up
respawn
exec /usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi-emperor.ini